# collector.cron
#!/bin/bash

/opt/kickgen/collector.py /opt/kickgen/all-rhel.txt 2>&1 > /dev/null

